# src/cron.d/chk-apex-domain
* * * * * root /usr/local/bin/start.sh chk-apex-domain >/dev/null 2>&1
