# src/cron.d/pubg-monitor
* * * * * root /usr/local/bin/start.sh pubg-monitor >/dev/null 2>&1
